# server/crontab
* * * * * curl http://spius_api:5000/crontab